# ctfs/SYD0x05/supervisord.conf
[supervisord]
nodaemon=true

[program:devcron]
command=/usr/local/bin/devcron.py /cron/crontab
